# example/agent-vm/hearth-agent-probe.service
[Unit]
Description=Hearth agent VM acceptance probe
After=multi-user.target systemd-networkd.service
Wants=systemd-networkd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hearth-agent-probe
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
