# example/hermes-vm/hermes.service
[Unit]
Description=Hermes Agent gateway
After=systemd-networkd.service
Wants=systemd-networkd.service

[Service]
Type=simple
User=agent
WorkingDirectory=/home/agent
Environment=HOME=/home/agent
Environment=PATH=/home/agent/.local/bin:/home/agent/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/agent/.local/bin/hermes serve --host 0.0.0.0 --port 9119
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
