# systemd/hearth-proxy.socket
[Unit]
Description=Agent-in-charge local Hearth socket proxy

[Socket]
ListenStream=/run/hearth.sock
SocketMode=0660

[Install]
WantedBy=sockets.target
