# systemd/hearth.service
[Unit]
Description=Hearth VM management daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
Environment=HEARTH_FIRMWARE=/var/lib/hearth/firmware/CLOUDHV.fd
Environment=PATH=/run/current-system/sw/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/hearthd
Restart=on-failure
RuntimeDirectory=hearth
RuntimeDirectoryMode=0770
StateDirectory=hearth
LogsDirectory=hearth
TimeoutStopSec=120s

[Install]
WantedBy=multi-user.target
