# packaging/systemd/hearth.service
[Unit]
Description=Hearth VM management daemon
Documentation=file:/usr/share/doc/hearth/operations.md
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
Environment=PATH=/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/hearthd --guest-kernel /usr/lib/hearth/kernel/vmlinux
Restart=on-failure
RestartSec=2
RuntimeDirectory=hearth
RuntimeDirectoryMode=0770
StateDirectory=hearth
LogsDirectory=hearth
TimeoutStopSec=120s

[Install]
WantedBy=multi-user.target
